# www-servers/pound/files/pound.init-2.0
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

extra_commands="configtest"

description_configtest="Run syntax tests for configuration files."

POUND="/usr/sbin/pound"
CONFIGFILE="/etc/pound.cfg"

depend() {
	need net
}

configtest() {
	ebegin "Checking ${SVCNAME} configuration"
	checkconfig
	eend $?
}

checkconfig() {
	${POUND} -c >/dev/null 2>&1
	ret=$?
	if [ $ret -ne 0 ]; then
		eerror "${SVCNAME} has detected an error in your setup:"
		${POUND} -c
	fi

	return $ret
}

start() {
	checkconfig || return 1

	ebegin "Starting pound"
	if [ ! -f "${CONFIGFILE}" ]; then
		eend 1 "configfile ${CONFIGFILE} not found."
	fi
	ulimit -n 32786
	start-stop-daemon --quiet --start --exec ${POUND} -- -f ${CONFIGFILE} -p /var/run/pound.pid
	eend $?
}

stop() {
	ebegin "Stopping pound"
	start-stop-daemon --quiet --stop --pidfile /var/run/pound.pid
	eend $?
}
