# app-emulation/wa-linux-agent/files/waagent.initd
#!/sbin/openrc-run

command=/usr/sbin/waagent
pidfile=/run/waagent.pid
command_args=-daemon
command_background=true
name="Azure Linux Agent"

depend()
{
	need localmount
	use logger network
	after bootmisc modules
}
