# app-misc/qualys-cloud-agent/files/qualys-cloud-agent.initd
#!/sbin/openrc-run

command="/usr/local/qualys/cloud-agent/bin/qualys-cloud-agent"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
retry="90"

depend() {
    need net
    after firewall
}
