# dev-db/percona-toolkit/files/pt-heartbeat.initd
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

command="/usr/bin/pt-heartbeat"
command_args="-D percona --update -h localhost --create-table"
command_background="true"
start_stop_daemon_args="-u nagios"
pidfile="/run/pt-heartbeat.pid"

depend() {
	user net
	after mysql
}
